# infra/onsite/systemd/alarm-cd-pull.service
# /etc/systemd/system/alarm-cd-pull.service
# Cf docs/CD_DESIGN.md §3 (Pull cote noeuds prod).
# Tourne via alarm-cd-pull.timer toutes les 5 min.
[Unit]
Description=Pull images Docker GHCR :stable (idempotent)
Documentation=https://github.com/claude-murgat/Alarm2.0/blob/master/docs/CD_DESIGN.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
# L'user `alarm` est dans le groupe docker (cf docs/PROVISIONING_ONSITE.md §14),
# donc peut parler au socket Docker sans sudo.
User=alarm
Group=alarm
# Cree /run/alarm-cd-pull/ chmod 755 owned alarm:alarm avant ExecStart, nettoye
# a l'arret. Sans ca, le lock /run/alarm-cd-pull/lock du script echoue avec
# "Permission denied" (/run est root-only). Pattern systemd standard pour
# services User=non-root qui ont besoin d'un repertoire ephemere sous /run.
RuntimeDirectory=alarm-cd-pull
# Cree /var/lib/alarm/cd/ owned alarm:alarm (et /var/lib/alarm/ parent si
# absent), persistent entre runs. Sans ca, le `mkdir -p /var/lib/alarm/cd` du
# script echoue avec "Permission denied" (/var/lib est root-only). Y stocke
# digests.tsv (cf scripts/alarm-cd-pull.sh STATE_DIR).
StateDirectory=alarm/cd
# Limite de temps : un pull qui prend > 5 min suggere un probleme reseau ou un
# blob trop gros (~600 Mo backend). On laisse 10 min de marge avant kill.
TimeoutStartSec=600
ExecStart=/usr/local/bin/alarm-cd-pull
# Pas de Restart= (Type=oneshot). Le timer relance toutes les 5 min de toute facon.
# Pas de RemainAfterExit (on veut un nouveau run a chaque tick, pas un service
# qui reste "active" entre les ticks).
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
